# ops/backup/systemd/dashboard-postgres-backup.timer
[Unit]
Description=Run Dashboard PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
RandomizedDelaySec=30m
Persistent=true
Unit=dashboard-postgres-backup.service

[Install]
WantedBy=timers.target
